# templates/default/nagios-nrpe-server.erb
#!/bin/sh
#
# Autogenerated by Chef.
#
# chkconfig: 2345 80 30
# description: nrpe is a daemon for a remote nagios server, \
#              running nagios plugins on this host.
# processname: <%= node['nrpe']['service_name'] %>
# config: <%= node['nrpe']['conf_dir'] %>/nrpe.cfg

### BEGIN INIT INFO
# Provides:          <%= node['nrpe']['service_name'] %>
# Required-Start:    $local_fs $remote_fs $syslog $named $network $time
# Required-Stop:     $local_fs $remote_fs $syslog $named $network
# Should-Start:      
# Should-Stop:       
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: nagios remote plugin executor
# Description:       nagios-nrpe is a daemon for running checks on remote systems.
### END INIT INFO

<% if ['rhel', 'fedora'].include?(node['platform_family']) -%>
# Source function library
if [ -f /etc/rc.d/init.d/functions ]; then
. /etc/rc.d/init.d/functions
elif [ -f /etc/init.d/functions ]; then
. /etc/init.d/functions
elif [ -f /etc/rc.d/functions ]; then
. /etc/rc.d/functions
fi


# Source networking configuration.
. /etc/sysconfig/network

# Check that networking is up.
[ ${NETWORKING} = "no" ] && exit 0

<% else -%>
. /lib/lsb/init-functions

<% end -%>
NrpeBin=<%= node['nrpe']['bin_dir'] %>/nrpe
NrpeCfg=<%= node['nrpe']['conf_dir'] %>/nrpe.cfg
PidFile=<%= node['nrpe']['pid_file'] %>

case "$1" in
start)
  # Start daemon.
  echo -n "Starting <%= node['nrpe']['service_name'] %>: "
  <% if ['rhel', 'fedora'].include?(node['platform_family']) -%>daemon $NrpeBin -c $NrpeCfg -d<% else -%>start-stop-daemon --start --quiet --exec $NrpeBin -- -c $NrpeCfg -d<% end -%>
  echo
  ;;
stop)
  # Stop daemon.
  echo -n "Shutting down <%= node['nrpe']['service_name'] %>: "
  <% if ['rhel', 'fedora'].include?(node['platform_family']) -%>killproc <%= node['nrpe']['service_name'] %><% else -%>start-stop-daemon --stop --quiet --exec $NrpeBin<% end -%>
  echo
  ;;
restart)
  $0 stop
  sleep 1
  $0 start
  ;;
reload)
  # reload the daemon.
  echo "Reloading <%= node['nrpe']['service_name'] %> configuration files"
  <% if ['rhel', 'fedora'].include?(node['platform_family']) -%>killproc -p $PidFile $NrpeBin -HUP<% else -%>start-stop-daemon --stop --signal HUP --quiet --pidfile $PidFile<% end -%>
  ;;
status)
  <% if ['rhel', 'fedora'].include?(node['platform_family']) -%>status<% else -%>status_of_proc<% end -%> -p $PidFile <%= node['nrpe']['service_name'] %> && exit 0 || exit $?
  ;;
*)
  echo "Usage: <%= node['nrpe']['service_name'] %> {start|stop|restart|status}"
  exit 1
esac

exit 0
