# templates/default/nrpe.service.erb
[Unit]
Description=Nagios Remote Program Executor
Documentation=http://www.nagios.org/documentation
Conflicts=nrpe.socket
Requires=network-online.target
After=var-run.mount nss-lookup.target network.target local-fs.target time-sync.target
Before=getty@tty1.service xdm.service

[Install]
WantedBy=multi-user.target

[Service]
Type=simple
User=<%= @nrpe['user'] %>
Group=<%= @nrpe['group'] %>
<% if ['rhel', 'fedora'].include?(node['platform_family']) && ::File.exist?('/etc/sysconfig/nrpe') -%>
EnvironmentFile=/etc/sysconfig/nrpe
<% end -%>
ExecStart=/usr/sbin/nrpe -c <%= @nrpe['conf_dir'] %>/nrpe.cfg -f $NRPE_SSL_OPT
ExecReload=/bin/kill -HUP $MAINPID
ExecStopPost=/bin/rm -f <%= @nrpe['pid_file'] %>
