# templates/default/nrpe.service.erb
[Unit]
Description=NRPE
After=network.target
Requires=network.target

[Service]
Type=forking
User=<%= @nrpe['user'] %>
Group=<%= @nrpe['group'] %>
<% if ['rhel', 'fedora'].include?(node['platform_family']) && ::File.exist?('/etc/sysconfig/nrpe') -%>
EnvironmentFile=/etc/sysconfig/nrpe
<% end -%>
ExecStart=/usr/sbin/nrpe -c <%= @nrpe['conf_dir'] %>/nrpe.cfg -d $NRPE_SSL_OPT
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
